# llm-bridge.service
[Unit]
Description=LLM Bridge — multi-harness session server
After=network.target

[Service]
Type=simple
User=__USER__
ExecStart=/usr/local/bin/llm-bridge
Environment=LLMBRIDGE_LISTEN_ADDR=:8160
Environment=LLMBRIDGE_IMAGES_DIR=__HOME__/repos/llm-bridge-server/images
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
